# systemd/machine-vision.service
[Unit]
Description=Machine Vision Pi Inspection Service
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/machine-vision-pi
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/machine-vision-pi/.venv/bin/python -m app.main
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
